# packaging/debian/postinst
#!/bin/sh
set -e

case "$1" in
	configure)
		if [ -d /run/systemd/system ]; then
			# Create the ssh-alertd system user and hand it the config so the
			# daemon can read its token (mode stays 0640 — not world-readable).
			systemd-sysusers /usr/lib/sysusers.d/ssh-alertd.conf >/dev/null 2>&1 || true
			systemd-tmpfiles --create /usr/lib/tmpfiles.d/ssh-alertd.conf >/dev/null 2>&1 || true
			# Belt-and-suspenders: ensure the daemon user can read its token.
			if getent passwd ssh-alertd >/dev/null 2>&1; then
				chown ssh-alertd:ssh-alertd /etc/ssh-alertd/config.json 2>/dev/null || true
				chmod 0640 /etc/ssh-alertd/config.json 2>/dev/null || true
			fi
			systemctl daemon-reload >/dev/null 2>&1 || true
		fi
		echo "ssh-alertd: edit /etc/ssh-alertd/config.json (set bot_token & chat_id),"
		echo "            then run: systemctl enable --now ssh-alertd"
		;;
esac

exit 0
