# deploy/ssh-alertd.service
[Unit]
Description=SSH Alert Daemon (sends an alert on successful SSH login)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ssh-alertd -config /etc/ssh-alertd/config.json
Restart=on-failure
RestartSec=5

# Runs as a dedicated system user (created via sysusers.d) that owns
# /etc/ssh-alertd/config.json (set up by tmpfiles.d) and belongs to
# systemd-journal so the journald log source works. For the "file" source, also
# grant read access to the auth log (e.g. add group "adm" on Debian/Ubuntu).
User=ssh-alertd
Group=ssh-alertd
SupplementaryGroups=systemd-journal
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
