# deploy/ssh-alertd.service
[Unit]
Description=SSH Alert Daemon (sends an alert on successful SSH login)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ssh-alertd -config /etc/ssh-alertd/config.json
Restart=on-failure
RestartSec=5

# Hardening. journald access needs the systemd-journal group; file mode needs
# read access to the auth log instead.
DynamicUser=yes
SupplementaryGroups=systemd-journal
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
